# ops/systemd/carbonet-auto-deploy.timer
[Unit]
Description=Check Carbonet main for deployable commits

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=15s
Persistent=true
Unit=carbonet-auto-deploy.service

[Install]
WantedBy=timers.target
