# ops/systemd/resonance-hermes-framework-qwen40-exl3.service.d/10-24453-firewall.conf
[Service]
ExecStartPost=/bin/bash -lc 'iptables -C DOCKER-USER -p tcp -s 172.16.1.231 -m conntrack --ctorigdstport 24453 -j ACCEPT 2>/dev/null || iptables -I DOCKER-USER 1 -p tcp -s 172.16.1.231 -m conntrack --ctorigdstport 24453 -j ACCEPT; iptables -C DOCKER-USER -p tcp -s 10.244.0.0/16 -m conntrack --ctorigdstport 24453 -j ACCEPT 2>/dev/null || iptables -I DOCKER-USER 2 -p tcp -s 10.244.0.0/16 -m conntrack --ctorigdstport 24453 -j ACCEPT; iptables -C DOCKER-USER -p tcp -m conntrack --ctorigdstport 24453 -j DROP 2>/dev/null || iptables -A DOCKER-USER -p tcp -m conntrack --ctorigdstport 24453 -j DROP'
